# .local/cfg/95-battery.rules
SUBSYSTEM=="power_supply", ACTION=="change", \
    RUN+="/usr/bin/bash -c '/usr/bin/polybar-msg hook battery 1 &>/dev/null'"
